# system/etc/profile.d/nix.sh
#!/bin/sh
if [ -e /home/cjr/.nix-profile/etc/profile.d/nix.sh ]; then . /home/cjr/.nix-profile/etc/profile.d/nix.sh; fi # added by Nix installer
