# system/usr/lib/systemd/system/daily@.timer
[Unit]
Description=Run %i every day

[Timer]
OnCalendar=daily
Persistent=true
Unit=%i.service

[Install]
WantedBy=timers.target
