# system/usr/lib/systemd/system/hourly@.timer
[Unit]
Description=Run %i every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=%i.service

[Install]
WantedBy=timers.target
